# bankserver.service
[Unit]
Description=Bank Server Service

[Service]
Type=simple
ExecStart=/bin/bash /home/ubuntu/bankserver

[Install]
WantedBy=multi-user.target